# ProjetECommerce/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/usr/sbin/nginx -g 'daemon off;'

[program:php-fpm]
logfile = /var/log/php-fpm/php-fpm.log
command= php-fpm -F -R
